# spec/fixtures/example/units/app.service
[Unit]
Name=app
Description=An app
ExecStart=/bin/app